# rss-notify.service
[Unit]
Description=Arch RSS Package Update Notifier
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/rss-notify
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
